# com.antigravity.rag.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.antigravity.rag</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>python3 "start_remote.py"</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/kaz/Library/Mobile Documents/com~apple~CloudDocs/antigravity/architectural_rag</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/kaz/Library/Mobile Documents/com~apple~CloudDocs/antigravity/architectural_rag/logs/launchd_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/kaz/Library/Mobile Documents/com~apple~CloudDocs/antigravity/architectural_rag/logs/launchd_stderr.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
